# pkgs/zsa-auto-us-layout/src/60-zsa-auto-us-layout.rules
SUBSYSTEMS=="usb", ATTRS{idVendor}=="3297", ATTRS{idProduct}=="4976",RUN+="@out@/bin/switch-layout"
